# market-monitor.service
[Unit]
Description=Market Monitor — Alertas de Câmbio e Commodities
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/market-monitor
ExecStart=/home/ubuntu/market-monitor/venv/bin/python monitor.py
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal
TimeoutStartSec=60

[Install]
WantedBy=multi-user.target
